# root/etc/services.d/sonarr/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

exec s6-setuidgid hotio "${APP_DIR}/bin/Sonarr" --nobrowser --data="${CONFIG_DIR}"
