# root/etc/s6-overlay/s6-rc.d/service-sonarr/run
#!/command/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

cd "${APP_DIR}/bin" || exit 1
exec s6-setuidgid hotio "${APP_DIR}/bin/Sonarr" --nobrowser --data="${CONFIG_DIR}"
